# profiles/ubuntu/rootfs_before/etc/xdg/autostart/fcitx5.desktop
[Desktop Entry]
Type=Application
Name=Fcitx5
Exec=fcitx5
X-GNOME-Autostart-enabled=true
NoDisplay=false
